# FirstYear/2nd/CSI/TS/TP3/scripts/monitord.timer
[Unit]
Description=Run monitord 1min timer
After=monitord-init.service
Requires=monitord-init.service

[Timer]
OnActiveSec=0s
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
