# FirstYear/2nd/CSI/TS/TP3/scripts/monitord.service
[Unit]
Description=Monitord

[Service]
ExecStart=/opt/monitord
